# cron/crontab
0 */6 * * * python /app/hello.py > /dev/stdout
